# tools/backup/systemd/asko-nightly-git-sync.timer
[Unit]
Description=Run Asko Nightly Git Sync

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
